# templates/cron.d/puppet_dns
* * * * * root /usr/local/bin/puppet_dns.rb
